# ansible/roles/web/r3status/templates/cache_realraum_imported_calendar.service.j2
[Unit]
Description=grical ics caching service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cache_realraum_imported_calendar.sh
User=spaceapi

[Install]
WantedBy=multi-user.target
