# ansible/roles/web/r3status/templates/cache_realraum_imported_calendar.timer.j2
[Unit]
Description=grical ics caching time

[Timer]
OnActiveSec=60
OnUnitActiveSec=60

[Install]
WantedBy=multi-user.target
